# roles/greenbone/templates/greenbone-admin.timer.j2
[Unit]
Description=Retry Greenbone administrator creation until gvmd is ready

[Timer]
OnActiveSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
