# roles/motd/templates/etc/update-motd.d/20-update.j2
#!/bin/bash
#
# Concise login MOTD. Only cheap, glanceable info shows here; the heavy
# collectors (open ports via ss, per-service systemctl, per-mount disk) live in
# the `sysinfo` command so they run on demand instead of on every SSH login.
# 16-color only and background-safe — readable on light and dark terminals.

# Colors
G='\e[{{ motd_color_primary }}m'      # frame + art (green)
W='\e[{{ motd_color_secondary }}m'    # emphasis (bold, no hardcoded fg)
CYN='\e[0;36m'
RED='\e[0;31m'
YEL='\e[0;33m'
R='\e[0m'

{% raw %}
HR="━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"

top()    { echo -e "${G}┏${HR}┓${R}"; }
bottom() { echo -e "${G}┗${HR}┛${R}"; }
blank()  { printf "${G}┃${R}%75s${G}┃${R}\n" ""; }

# Pad using the BYTE length of plain text. Values are ASCII only, so alignment
# holds; multibyte chars live only in the frame, outside the %s.
pline() {
  local label="$1" value="$2" color="${3:-$W}"
  local visible="${label}${value}"
  local padding=$((74 - ${#visible}))
  [ $padding -lt 0 ] && padding=0
  printf "${G}┃${R} %s${color}%s${R}%${padding}s${G}┃${R}\n" "$label" "$value" ""
}

# --- Gather (cheap, login-time) ---
UPTIME=$(uptime -p 2>/dev/null | sed 's/^up //' || uptime)
read -r LOAD1 LOAD5 LOAD15 _ < /proc/loadavg
CPUS=$(nproc 2>/dev/null)
HNAME=$(hostname -s 2>/dev/null)
DISTRO=$(. /etc/os-release 2>/dev/null && echo "$NAME $VERSION_ID")
RAM_TOTAL=$(free -h | awk '/^Mem:/{print $2}')
RAM_USED=$(free -h | awk '/^Mem:/{print $3}')
RAM_PCT=$(free | awk '/^Mem:/{printf "%.0f", $3/$2*100}')
SWAP_TOTAL=$(free -h | awk '/^Swap:/{print $2}')
SWAP_USED=$(free -h | awk '/^Swap:/{print $3}')
IP=$(ip -o -4 addr show scope global 2>/dev/null | awk '{print $4}' | cut -d/ -f1 | paste -sd' ')
GW=$(ip route 2>/dev/null | awk '/^default/{print $3; exit}')
DNS=$(resolvectl status 2>/dev/null | grep -A1 'DNS Servers' | tail -1 | sed 's/^ *//;s/^DNS Servers: //')
[ -z "$DNS" ] && DNS=$(awk '/^nameserver/{printf "%s ", $2}' /etc/resolv.conf 2>/dev/null)
[ -z "$DNS" ] && DNS="N/A"
DISK=$(df -h -x tmpfs -x devtmpfs -x squashfs 2>/dev/null | awk 'NR>1 && $6 !~ /boot/{printf "%s %s  ", $6, $5}' | sed 's/  *$//')

RAM_CLR=$G
[ "$RAM_PCT" -gt 75 ] 2>/dev/null && RAM_CLR=$YEL
[ "$RAM_PCT" -gt 85 ] 2>/dev/null && RAM_CLR=$RED

# --- Concise summary ---
top
pline "Host   " "$HNAME | $DISTRO | ${CPUS} CPU"
pline "Up     " "$UPTIME | load $LOAD1 $LOAD5 $LOAD15"
pline "RAM    " "$RAM_USED / $RAM_TOTAL (${RAM_PCT}%) | swap $SWAP_USED / $SWAP_TOTAL" "$RAM_CLR"
pline "Disk   " "$DISK"
pline "Net    " "$IP | gw $GW"
pline "DNS    " "$DNS"
bottom
{% endraw %}

# --- Branding banner ---
top
blank
{{ motd_banner_managed_line }}
blank
{{ motd_banner_art }}
{{ motd_banner_tagline }}
blank
bottom
{% if motd_show_support | default(false) %}

# --- Support (optional: set motd_show_support: true) ---
top
pline "Support:  " "{{ motd_sysadmins_signature }}"
pline "Email:    " "{{ motd_sysadmins_email }}" "$CYN"
pline "Info:     " "https://vagrant-zones.startcloud.com/changes" "$CYN"
pline "Source:   " "https://boxvault.startcloud.com/{{ vendor | default('STARTcloud') }}/{{ distro }}{{ distro_version }}-{{ vm_type }}" "$CYN"
bottom
{% endif %}

echo -e "  Run ${W}sysinfo${R} for open ports, services and full system detail."
